# client/klyra.service
[Unit]
Description=Klyra AI Companion
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/klyramachine/client
ExecStart=/usr/bin/python3 /home/pi/klyramachine/client/client_companion.py
Restart=always
RestartSec=10

# Environment variables (if needed)
# Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
